# recharge/execucao_mensal.cron
00 12 10 * * /projects/phone-recharge/run.sh >> /projects/phone-recharge/logs/recharge.log 2>&1